# pi/mavlink.service
[Unit]
Description=Weed pilot (STM32 serial <-> GCS UDP, weed target relay)
After=network-online.target
Wants=network-online.target

[Service]
User=pi
WorkingDirectory=/home/pi/drone
EnvironmentFile=/home/pi/drone/.env
ExecStart=/usr/bin/python3 /home/pi/drone/weed_pilot.py
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
